# config/supervisord.conf
[supervisord]
nodaemon=true

[program:railo]
command=/bin/bash -c "exec /opt/railo/railo_ctl start -D FOREGROUND"

